# trias.service
[Unit]
Description=TRIAS API Server - Public Transit Information
After=network.target

[Service]
Type=exec
User=trias
Group=trias
WorkingDirectory=/opt/trias
Environment="PATH=/opt/trias/venv/bin"
# Set to true to use port 80 (requires capability), false for port 8080
Environment="USE_PORT_80=false"
ExecStart=/opt/trias/venv/bin/gunicorn --config /opt/trias/gunicorn_config.py app:app
# Uncomment the next line to allow binding to port 80 as non-root user
# AmbientCapabilities=CAP_NET_BIND_SERVICE
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
